# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
loglevel = INFO

[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:secret]
command=/bin/sh -c "[ ! -f './config/.secret_key' ] && SECRET=$(python3 -c 'import secrets; print(secrets.token_hex())') && echo SECRET_KEY=$SECRET > './config/.secret_key'"
directory=/app
user=root
autostart=true
priority=1

[program:alembic]
command=alembic upgrade head
user=root
directory=/app
autostart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2

[program:fastapi]
command=fastapi run app/main.py --port 8000
user=root
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3